# scripts/com.flowboard.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  TEMPLATE — paths marked __VENV_PATH__ and __WORKDIR__ are substituted by
  `bin/flowboard install` (see generate_plist() inside that script). The
  per-user copy in ~/Library/LaunchAgents/com.flowboard.agent.plist gets
  regenerated on every install; do not edit it by hand.
-->
<dict>
    <key>Label</key>
    <string>com.flowboard.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PATH__/bin/uvicorn</string>
        <string>flowboard.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8101</string>
        <string>--log-level</string>
        <string>info</string>
    </array>

    <!-- Run from agent/ so the `flowboard` package is importable. config.py
         resolves ROOT from __file__ and defaults storage/ to <project>/storage,
         so the SQLite file stays where the rest of the project expects it. -->
    <key>WorkingDirectory</key>
    <string>__WORKDIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__VENV_PATH__/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Start on launchctl load (== at user login). -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on any exit, including crash. ThrottleInterval=5s avoids a
         crash-loop burning CPU if the binary is missing or DB is locked. -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>

    <!-- Logs: split stdout/stderr so HTTP 4xx doesn't drown real errors. -->
    <key>StandardOutPath</key>
    <string>/Users/friday/.flowboard/logs/agent.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/friday/.flowboard/logs/agent.err.log</string>

    <!-- ProcessType=Interactive keeps the service in the user's GUI session
         (some macOS subsystems behave differently under Background). -->
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
